# systemd/minecraft/minecraft-restart.service
[Unit]
Description=Safe Restart Minecraft Server

[Service]
Type=oneshot
ExecStart=/usr/local/bin/minecraft-safe-restart.sh