# systemd/minecraft/minecraft-backup.service
[Unit]
Description=Minecraft Backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/minecraft-backup.sh